# xhs-auto-publisher/deploy/systemd/xhs-auto-publisher-cloud.service
[Unit]
Description=XHS Auto Publisher Cloud
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=%h/projects/xhs-auto-publisher
EnvironmentFile=-%h/projects/xhs-auto-publisher/.env
ExecStart=%h/projects/xhs-auto-publisher/deploy/run_with_xvfb.sh %h/projects/xhs-auto-publisher/examples/openclaw_business_content.json
Restart=no
TimeoutStartSec=900
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
